# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true

[program:php-fpm]
command=php-fpm8 --nodaemonize --fpm-config /etc/php8/php-fpm.conf
autorestart=true
